# crontab
* * * * * root /putio_script.sh > /putio.log
# end
